# scripts/cloudflared/com.cloudflare.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cloudflare.tunnel</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/cloudflared</string>
    <string>tunnel</string>
    <string>run</string>
    <string>alpaca-cam</string>
  </array>

  <key>KeepAlive</key>
  <true/>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/alpaca/cloudflared/tunnel.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/alpaca/cloudflared/tunnel.err</string>

  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
